# systemd/valentina-restore.timer
[Unit]
Description=Valentina Self-Restoration timer — 15s after boot
[Timer]
OnBootSec=15s
Persistent=true
[Install]
WantedBy=timers.target
